# conf/init/couchdb.conf
description "CouchDB server"
author "Dustin Moore"

start on startup
stop on shutdown

script
  export HOME="/root"
  exec /usr/local/bin/couchdb
end script

respawn
