# conf/init/http_server.conf
description "http NodeJS server"
author "Dustin Moore"

start on startup
stop on shutdown

script
  export HOME="/root"
  exec /usr/local/bin/node /home/dustin/PicList/server.js
end script

respawn
